# roles/gui/files/bg.service
[Unit]
Description=Rotate wallpaper every 30 minutes

[Service]
ExecStart=/usr/bin/feh --randomize --bg-scale /home/gavin/.wallpapers/
